# systemd/dgx-moa-resident.target
[Unit]
Description=DGX MoA resident profile
Requires=dgx-moa-gateway.service dgx-moa-executor.service dgx-moa-planner.service dgx-moa-reviewer.service
After=dgx-moa-gateway.service dgx-moa-executor.service dgx-moa-planner.service dgx-moa-reviewer.service dgx-moa-judge.target
Conflicts=dgx-moa-judge.target

[Install]
WantedBy=default.target

